# ambalinux_sdk/ambarella/package/ambarella_customize/source/etc_sysvinit/inittab
# /etc/inittab
#
# This inittab is a basic inittab sample for sysvinit, which mimics
# Buildroot's default inittab for BusyBox.
id:3:initdefault:

si0::sysinit:/bin/mount -t proc proc /proc
#si1::sysinit:/bin/mount -o remount,rw /
si2::sysinit:/bin/mkdir -p /dev/pts
si3::sysinit:/bin/mkdir -p /dev/shm
si4::sysinit:/bin/mount -a
si5::sysinit:/bin/hostname -F /etc/hostname
rcS:12345:wait:/etc/init.d/rcS

# Put a getty on the serial port
S1::respawn:/bin/sh

# ttyUSB0 on usb host.
# ttyUSB0::respawn:-/sbin/getty  115200 ttyUSB0

# ttyGS0 on usb device
# ttyGS0::respawn:-/bin/sh

# Stuff to do for the 3-finger salute
#::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
shd0:06:wait:/etc/init.d/rcK
shd1:06:wait:/sbin/swapoff -a
shd2:06:wait:/bin/umount -a -r

# The usual halt or reboot actions
#hlt0:0:wait:/sbin/halt -dhp
#reb0:6:wait:/sbin/reboot
