# ambalinux_sdk/ambarella/bootloader/boards/cv2fsbub/rootfs/default/etc/udev/rules.d/09-load-ucode.rules
# Load ucode when IAV driver is loaded

SUBSYSTEM=="ambpriv", KERNEL=="iav_v5", RUN+="/usr/sbin/load_ucode /lib/firmware"
