# ambalinux_sdk/ambarella/bootloader/boards/cv2bub/rootfs/default/etc/udev/rules.d/09-load-ucode.rules
# Load ucode when IAV driver is loaded

SUBSYSTEM=="ambpriv", KERNEL=="iav_s5l", RUN+="/usr/sbin/load_ucode /lib/firmware"
